# deployment/room-booking-expiry.timer
[Unit]
Description=Run room booking expiry every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=5s
Persistent=true
Unit=room-booking-expiry.service

[Install]
WantedBy=timers.target
